# dotfiles/systemd/user/bar-htop.service
[Unit]
Description=htop hidden window for status bar

[Service]
ExecStart=kitty --name=bar-htop htop
Restart=always

[Install]
WantedBy=i3.target
