# dotfiles/systemd/user/bitwarden-desktop.service
[Unit]
Description=Bitwarden Desktop

[Service]
ExecStart=bitwarden-desktop

[Install]
WantedBy=xsession.target
