# dotfiles/systemd/user/dunst.service
[Unit]
Description=Dunst notification manager
Type=dbus

[Service]
ExecStart=dunst

[Install]
WantedBy=xsession.target
