# dotfiles/systemd/user/i3.target
[Unit]
Description=i3 is running
BindsTo=graphical-session.target
