# dotfiles/systemd/user/nm-applet.service
[Unit]
Description=Network manager applet

[Service]
ExecStart=nm-applet

[Install]
WantedBy=xsession.target
