# dotfiles/systemd/user/picom.service
[Unit]
Description=Picom compositor

[Service]
ExecStart=picom

[Install]
WantedBy=xsession.target
