# dotfiles/systemd/user/udiskie.service
[Unit]
Description=udiskie automounter
Type=dbus

[Service]
ExecStart=udiskie

[Install]
WantedBy=xsession.target
