# dotfiles/systemd/user/unclutter.service
[Unit]
Description=Unclutter

[Service]
ExecStart=unclutter --jitter 20

[Install]
WantedBy=xsession.target
