# dotfiles/systemd/user/xsession.target
[Unit]
Description=X session is running
BindsTo=graphical-session.target
